# kaikei-bot.service
[Unit]
Description=Discord 会計申請ボット
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=haruka
WorkingDirectory=/home/haruka/ドキュメント/kaikei
ExecStart=/home/haruka/ドキュメント/kaikei/.venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
